# 99-camera-sync.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="$IDVENDOR", ATTR{idProduct}=="$IDPRODUCT", ACTION=="add", RUN+="/home/$USER/camera/sync_camera.sh"
